# systemd/user-services/set-random-theme.timer
[Unit]
Description=Change theme periodically

[Timer]
OnCalendar=0/1:00:00

[Install]
WantedBy=timers.target
