# systemd/user-services/syncthing.service
[Unit]
Description=syncthing service

[Service]
Type=simple
ExecStart=syncthing
Restart=always    
RestartSec=5

[Install]
WantedBy=default.target
