# systemd/user-services/volume-icon-applet.service
[Unit]
Description=volumeicon service

[Service]
Type=simple
ExecStart=volumeicon
Restart=always    
RestartSec=5

[Install]
WantedBy=default.target
