# systemd/nm-applet.service
[Unit]
Description=nm-applet service

[Service]
Type=simple
ExecStart=nm-applet
Restart=always    
RestartSec=5
