# i3/theme-changer/set-random-theme.timer
[Unit]
Description=Change theme periodically (every 2 hours)

[Timer]
OnCalendar=0/2:00:00
